# autostart/blueman-manager.desktop
[Desktop Entry]
Name=Bluetooth Manager
Type=Application
TryExec=blueman-applet
Exec=blueman-applet
