# autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Type=Application
TryExec=dunst
Exec=dunst
