# autostart/nitrogen.desktop
[Desktop Entry]
Name=nitrogen
Type=Application
TryExec=nitrogen --restore
Exec=nitrogen --restore
