# autostart/nm-applet.desktop
[Desktop Entry]
Name=NetworkManager Applet
Type=Application
TryExec=nm-applet
Exec=nm-applet
