# autostart/picom.desktop
[Desktop Entry]
Name=picom
Type=Application
TryExec=picom
Exec=picom
