# autostart/pipewire.desktop
[Desktop Entry]
Name=PipeWire
Type=Application
TryExec=pipewire
Exec=pipewire
