# autostart/xfce4-clipman.desktop
[Desktop Entry]
Name=Clipman
Type=Application
TryExec=xfce4-clipman
Exec=xfce4-clipman
